# supervisord.ini
[program:snapfile]
command=/home/yxr/.pyenv/versions/3.8.1/bin/snapfile
user=yxr
autostart=true
autorestart=true
startretries=3
exitcodes=0
redirect_stderr=true
stdout_logfile=/var/www/snapfile/logs/snapfile.out
