# supervisord.ini
[program:snapfile]
command=/home/yxr/.pyenv/versions/3.8.1/bin/snapfile
environment=ENV="PROD"
user=yxr
autostart=false
autorestart=true
startretries=3
exitcodes=0
redirect_stderr=true
stdout_logfile=/var/www/snapfile/logs/snapfile.log
